# global/profile.d/xdgninja.sh
export CARGO_HOME="$XDG_DATA_HOME/cargo"
export RUSTUP_HOME="$XDG_DATA_HOME/rustup"

export CUDA_CACHE_PATH="XDG_CACHE_HOME"/nv

#export GNUPGHOME="$XDG_DATA_HOME/gnupg"

export DOTNET_CLI_HOME="$XDG_DATA_HOME/dotnet"
export NUGET_PACKAGES="$XDG_CACHE_HOME/NuGetPackages"

export _JAVA_OPTIONS="-Djava.util.prefs.userRoot=$XDG_CONFIG_HOME/java -Djavafx.cachedir=$XDG_CACHE_HOME/openjfx"
export GRADLE_USER_HOME="$XDG_DATA_HOME/gradle"

export PYTHONSTARTUP="$XDG_CONFIG_HOME/python/pythonrc"
export GTK2_RC_FILES="$XDG_CONFIG_HOME/gtk-2.0/gtkrc"

export KDEHOME="$XDG_CONFIG_HOME/kde"

export WINEPREFIX="$XDG_DATA_HOME/wine"

export HISTFILE="$XDG_STATE_HOME/zsh/history"
export SQLITE_HISTORY="$XDG_CACHE_HOME/sqlite_history"
export REDISCLI_HISTFILE="$XDG_DATA_HOME/redis/rediscli_history"
export ERRFILE="$XDG_CACHE_HOME/X11/xsession-errors"
export NODE_REPL_HISTORY="$XDG_DATA_HOME/node_repl_history"

export ARDUINO_DIRECTORIES_DATA="$XDG_CONFIG_HOME/arduino15"
export ARDUINO_sketchbook_location="$HOME/workspaces/arduino"
export ARDUINO_libraries_location="$ARDUINO_DIRECTORIES_DATA/libraries"

export DOCKER_CONFIG="$XDG_CONFIG_HOME/docker"

export XCURSOR_PATH=/usr/share/icons:$XDG_DATA_HOME/icons

export TEXMFVAR="$XDG_CACHE_HOME/texlive/texmf-var"

export W3M_DIR="$XDG_DATA_HOME/w3m"

#export XAUTHORITY="$XDG_RUNTIME_DIR/Xauthority"
